# packaging/vtesserad.service
[Unit]
Description=Vtessera metering daemon
Documentation=https://github.com/douglasdemaio/vtessera
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/vtesserad --config /etc/vtessera/vtessera.toml
DynamicUser=yes
StateDirectory=vtessera
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged ~@resources
CapabilityBoundingSet=
AmbientCapabilities=
# Default build opens no sockets:
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
# If built with --features submit, set instead:
#   RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
#   IPAddressDeny=any
#   IPAddressAllow=<your submission endpoint CIDR>
#
# Note: DO NOT set ProcSubset=pid here. The daemon must read /proc/meminfo,
# /proc/stat etc., which ProcSubset=pid hides.
#
# Note: DynamicUser=yes randomizes the runtime username. AppArmor profiles
# pathed by username will not match. The shipped profile uses peer/object
# rules, not user-pathed rules — keep it that way.
Restart=on-failure

[Install]
WantedBy=multi-user.target
